# root/usr/lib/systemd/system/init.service
[Unit]
Description=System init scipt from user

[Service]
Type=simple
ExecStart=/opt/init

[Install]
WantedBy=default.target
